# titan-api-v13.service
[Unit]
Description=Titan VMOS API Server (FastAPI + uvicorn)
Documentation=https://github.com/malithwishwa02-dot/titan-vmos
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/titan-v13-device/vmos-titan
Environment=PYTHONPATH=/opt/titan-v13-device/vmos-titan/vmos_titan/api:/opt/titan-v13-device/vmos-titan/vmos_titan/core:/opt/titan-v13-device/vmos-titan
EnvironmentFile=-/opt/titan-v13-device/.env
ExecStart=/opt/titan-v13-device/venv/bin/uvicorn vmos_titan.api.titan_api:app \
    --host 127.0.0.1 --port 8080 \
    --workers 2 --loop uvloop --http httptools
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=15

# Resource limits
MemoryMax=2G
CPUQuota=200%

# Security hardening
ProtectSystem=strict
ReadWritePaths=/opt/titan-v13-device/data /opt/titan-v13-device/venv /opt/titan/cuttlefish /tmp
ReadWritePaths=/opt/titan-v13-device/.env /opt/titan-v13-device/docker/ssl /opt/titan/data
ProtectHome=read-only
NoNewPrivileges=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target
